# root/etc/services.d/cron/run
#!/usr/bin/with-contenv bash

# update trackers regularly
exec crond -f -L /dev/null
